# com.hermes.r2d2-llm.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.hermes.r2d2-llm</string>
    <key>ProgramArguments</key>
    <array>
        <string>/usr/bin/python3</string>
        <string>/Users/iukha/r2d2_llm.py</string>
        <string>--api</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>KeepAlive</key>
    <true/>
    <key>StandardOutPath</key>
    <string>/Users/iukha/.hermes/r2d2.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/iukha/.hermes/r2d2.err</string>
    <key>WorkingDirectory</key>
    <string>/Users/iukha</string>
</dict>
</plist>
